# files/etc/systemd/system/dhcpcd.service.d/no-wait.conf
[Service]
Type=simple
ExecStart=
ExecStart=/usr/bin/dhcpcd -q -N -B
